# systemd/crawlify-update.timer
[Unit]
Description=Run Crawlify Kleine Anfragen Update daily
Documentation=https://github.com/rstemmler/crawlify-kleine-anfragen

[Timer]
# Run daily at 3:00 AM
OnCalendar=*-*-* 03:00:00
# Run immediately if last run was missed (e.g., system was off)
Persistent=true
# Random delay up to 15 minutes to avoid thundering herd
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
